# provision/roles/docker-registry/templates/publisher.service.j2
[Unit]
Description=Minios S3 to CernVM-FS publisher service

[Service]
ExecStart={{ workdir }}/bin/publisher
StandardOutput=journal
StandardError=journal
